# runComputeController
### BEGIN INIT INFO
# Provides:          runSmashComputeController
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Enable service provided by daemon.  This runs basic setup like bundle install
### END INIT INFO
#!/bin/sh
#!/usr/bin/env ruby

cd /home/ubuntu/smash_compute_cloud
bundle install

/home/ubuntu/.rbenv/shims/ruby compute_controller.rb &> computeController.log

exit 0